# docker/nginx/crontab

0 0 * * * /docker/log-backup.sh /nginx-log nginx_error
0 0 * * * /docker/log-backup.sh /nginx-log nginx_access
